# templates/default/shipper.upstart.erb
description "Start shipper"

start on startup
stop on shutdown

respawn
respawn limit 10 5

setuid <%= @app_user %>

env GH_KEY=\"<%= @github_key %>\"

exec \
<%= @shipper_path %>/shipper run -c <%= @config_path %> >> <%= @app_path %>/shipper.log 2>&1
